# sysdrv/tools/board/base_rootfs/etc/profile.d/RkEnv.sh
export OEM=/oem
export PATH=$PATH:$OEM:$OEM/bin:$OEM/usr/bin:$OEM/sbin:$OEM/usr/sbin
export LD_LIBRARY_PATH=$OEM/usr/lib:$OEM/lib:$LD_LIBRARY_PATH
